# dockerfiles/code/code/services.d/dind/run
#!/usr/bin/execlineb -P

rm -f /var/run/docker.pid
dind dockerd